# packages/systemd/FOGImageReplicator.service
[Unit]
Description=FOGImageReplicator
After=syslog.target network.target mysql.service

[Service]
PrivateTmp=true
ExecStart=/opt/fog/service/FOGImageReplicator/FOGImageReplicator &
ExecReload=/bin/kill -HUP $MAINPID
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
