# packages/systemd/FOGMulticastManager.service
[Unit]
Description=FOGMulticastManager
After=syslog.target network.target mysql.service

[Service]
PrivateTmp=true
ExecStart=/opt/fog/service/FOGMulticastManager/FOGMulticastManager &
ExecReload=/bin/kill -HUP $MAINPID
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
