# packages/systemd/FOGPingHosts.service
[Unit]
Description=FOGPingHosts
After=syslog.target network.target mysql.service

[Service]
PrivateTmp=true
ExecStart=/opt/fog/service/FOGPingHosts/FOGPingHosts &
ExecReload=/bin/kill -HUP $MAINPID
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
