# packages/systemd/FOGScheduler.service
[Unit]
Description=FOGScheduler
After=syslog.target network.target mysql.service

[Service]
PrivateTmp=true
ExecStart=/opt/fog/service/FOGTaskScheduler/FOGTaskScheduler &
ExecReload=/bin/kill -HUP $MAINPID
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
